# evidence/T1053.005/task_suspicious_executable.json
{
  "agent": {
    "ip": "x.x.x.x",
    "name": "windows-detection-vm",
    "id": "001"
  },
  "manager": {
    "name": "wazuh-manager"
  },
  "data": {
    "win": {
      "eventdata": {
        "subjectLogonId": "0x32b7647",
        "clientProcessId": "5376",
        "subjectUserSid": "S-1-5-21-1936990016-347023981-4024709258-1001",
        "taskContent": "&lt;?xml version=\\\"1.0\\\" encoding=\\\"UTF-16\\\"?&gt;  &lt;Task version=\\\"1.2\\\" xmlns=\\\"http://schemas.microsoft.com/windows/2004/02/mit/task\\\"&gt;    &lt;RegistrationInfo&gt;      &lt;Date&gt;2026-07-07T23:33:53&lt;/Date&gt;      &lt;Author&gt;hostname\\\\username&lt;/Author&gt;      &lt;URI&gt;\\\\T1053_005_OnStartup&lt;/URI&gt;    &lt;/RegistrationInfo&gt;    &lt;Triggers&gt;      &lt;BootTrigger&gt;        &lt;StartBoundary&gt;2026-07-07T23:33:00&lt;/StartBoundary&gt;        &lt;Enabled&gt;true&lt;/Enabled&gt;      &lt;/BootTrigger&gt;    &lt;/Triggers&gt;    &lt;Settings&gt;      &lt;MultipleInstancesPolicy&gt;IgnoreNew&lt;/MultipleInstancesPolicy&gt;      &lt;DisallowStartIfOnBatteries&gt;true&lt;/DisallowStartIfOnBatteries&gt;      &lt;StopIfGoingOnBatteries&gt;true&lt;/StopIfGoingOnBatteries&gt;      &lt;AllowHardTerminate&gt;true&lt;/AllowHardTerminate&gt;      &lt;StartWhenAvailable&gt;false&lt;/StartWhenAvailable&gt;      &lt;RunOnlyIfNetworkAvailable&gt;false&lt;/RunOnlyIfNetworkAvailable&gt;      &lt;IdleSettings&gt;        &lt;Duration&gt;PT10M&lt;/Duration&gt;        &lt;WaitTimeout&gt;PT1H&lt;/WaitTimeout&gt;        &lt;StopOnIdleEnd&gt;true&lt;/StopOnIdleEnd&gt;        &lt;RestartOnIdle&gt;false&lt;/RestartOnIdle&gt;      &lt;/IdleSettings&gt;      &lt;AllowStartOnDemand&gt;true&lt;/AllowStartOnDemand&gt;      &lt;Enabled&gt;true&lt;/Enabled&gt;      &lt;Hidden&gt;false&lt;/Hidden&gt;      &lt;RunOnlyIfIdle&gt;false&lt;/RunOnlyIfIdle&gt;      &lt;WakeToRun&gt;false&lt;/WakeToRun&gt;      &lt;ExecutionTimeLimit&gt;PT72H&lt;/ExecutionTimeLimit&gt;      &lt;Priority&gt;7&lt;/Priority&gt;    &lt;/Settings&gt;    &lt;Actions Context=\\\"Author\\\"&gt;      &lt;Exec&gt;        &lt;Command&gt;cmd.exe&lt;/Command&gt;        &lt;Arguments&gt;/c calc.exe&lt;/Arguments&gt;      &lt;/Exec&gt;    &lt;/Actions&gt;    &lt;Principals&gt;      &lt;Principal id=\\\"Author\\\"&gt;        &lt;UserId&gt;S-1-5-18&lt;/UserId&gt;        &lt;RunLevel&gt;LeastPrivilege&lt;/RunLevel&gt;      &lt;/Principal&gt;    &lt;/Principals&gt;  &lt;/Task&gt;",
        "parentProcessId": "6080",
        "rpcCallClientLocality": "0",
        "subjectDomainName": "hostname",
        "clientProcessStartKey": "9007199254743556",
        "taskName": "\\\\T1053_005_OnStartup",
        "fQDN": "hostname",
        "subjectUserName": "username"
      },
      "system": {
        "eventID": "4698",
        "keywords": "0x8020000000000000",
        "providerGuid": "{54849625-5478-4994-a5ba-3e3b0328c30d}",
        "level": "0",
        "channel": "Security",
        "opcode": "0",
        "message": "\"A scheduled task was created.\r\n\r\nSubject:\r\n\tSecurity ID:\t\tS-1-5-21-1936990016-347023981-4024709258-1001\r\n\tAccount Name:\t\tusername\r\n\tAccount Domain:\t\thostname\r\n\tLogon ID:\t\t0x32B7647\r\n\r\nTask Information:\r\n\tTask Name: \t\t\\T1053_005_OnStartup\r\n\tTask Content: \t\t<?xml version=\"1.0\" encoding=\"UTF-16\"?>\r\n<Task version=\"1.2\" xmlns=\"http://schemas.microsoft.com/windows/2004/02/mit/task\">\r\n  <RegistrationInfo>\r\n    <Date>2026-07-07T23:33:53</Date>\r\n    <Author>hostname\\username</Author>\r\n    <URI>\\T1053_005_OnStartup</URI>\r\n  </RegistrationInfo>\r\n  <Triggers>\r\n    <BootTrigger>\r\n      <StartBoundary>2026-07-07T23:33:00</StartBoundary>\r\n      <Enabled>true</Enabled>\r\n    </BootTrigger>\r\n  </Triggers>\r\n  <Settings>\r\n    <MultipleInstancesPolicy>IgnoreNew</MultipleInstancesPolicy>\r\n    <DisallowStartIfOnBatteries>true</DisallowStartIfOnBatteries>\r\n    <StopIfGoingOnBatteries>true</StopIfGoingOnBatteries>\r\n    <AllowHardTerminate>true</AllowHardTerminate>\r\n    <StartWhenAvailable>false</StartWhenAvailable>\r\n    <RunOnlyIfNetworkAvailable>false</RunOnlyIfNetworkAvailable>\r\n    <IdleSettings>\r\n      <Duration>PT10M</Duration>\r\n      <WaitTimeout>PT1H</WaitTimeout>\r\n      <StopOnIdleEnd>true</StopOnIdleEnd>\r\n      <RestartOnIdle>false</RestartOnIdle>\r\n    </IdleSettings>\r\n    <AllowStartOnDemand>true</AllowStartOnDemand>\r\n    <Enabled>true</Enabled>\r\n    <Hidden>false</Hidden>\r\n    <RunOnlyIfIdle>false</RunOnlyIfIdle>\r\n    <WakeToRun>false</WakeToRun>\r\n    <ExecutionTimeLimit>PT72H</ExecutionTimeLimit>\r\n    <Priority>7</Priority>\r\n  </Settings>\r\n  <Actions Context=\"Author\">\r\n    <Exec>\r\n      <Command>cmd.exe</Command>\r\n      <Arguments>/c calc.exe</Arguments>\r\n    </Exec>\r\n  </Actions>\r\n  <Principals>\r\n    <Principal id=\"Author\">\r\n      <UserId>S-1-5-18</UserId>\r\n      <RunLevel>LeastPrivilege</RunLevel>\r\n    </Principal>\r\n  </Principals>\r\n</Task>\r\n\r\nOther Information:\r\n\tProcessCreationTime: \t\t9007199254743556\r\n\tClientProcessId: \t\t\t5376\r\n\tParentProcessId: \t\t\t6080\r\n\tFQDN: \t\t0\r\n\t\"",
        "version": "1",
        "systemTime": "2026-07-08T06:33:53.3441984Z",
        "eventRecordID": "56676",
        "threadID": "2536",
        "computer": "hostname",
        "task": "12804",
        "processID": "972",
        "severityValue": "AUDIT_SUCCESS",
        "providerName": "Microsoft-Windows-Security-Auditing"
      }
    }
  },
  "rule": {
    "firedtimes": 2,
    "mail": true,
    "level": 15,
    "description": "High Risk: Scheduled task running suspicious executable: \\\\T1053_005_OnStartup",
    "groups": [
      "windows",
      "persistence",
      "sigma_coverage"
    ],
    "mitre": {
      "technique": [
        "Scheduled Task"
      ],
      "id": [
        "T1053.005"
      ],
      "tactic": [
        "Execution",
        "Persistence",
        "Privilege Escalation"
      ]
    },
    "id": "100006"
  },
  "decoder": {
    "name": "windows_eventchannel"
  },
  "full_log": "{\"win\":{\"system\":{\"providerName\":\"Microsoft-Windows-Security-Auditing\",\"providerGuid\":\"{54849625-5478-4994-a5ba-3e3b0328c30d}\",\"eventID\":\"4698\",\"version\":\"1\",\"level\":\"0\",\"task\":\"12804\",\"opcode\":\"0\",\"keywords\":\"0x8020000000000000\",\"systemTime\":\"2026-07-08T06:33:53.3441984Z\",\"eventRecordID\":\"56676\",\"processID\":\"972\",\"threadID\":\"2536\",\"channel\":\"Security\",\"computer\":\"hostname\",\"severityValue\":\"AUDIT_SUCCESS\",\"message\":\"\\\"A scheduled task was created.\\r\\n\\r\\nSubject:\\r\\n\\tSecurity ID:\\t\\tS-1-5-21-1936990016-347023981-4024709258-1001\\r\\n\\tAccount Name:\\t\\tusername\\r\\n\\tAccount Domain:\\t\\thostname\\r\\n\\tLogon ID:\\t\\t0x32B7647\\r\\n\\r\\nTask Information:\\r\\n\\tTask Name: \\t\\t\\\\T1053_005_OnStartup\\r\\n\\tTask Content: \\t\\t<?xml version=\\\"1.0\\\" encoding=\\\"UTF-16\\\"?>\\r\\n<Task version=\\\"1.2\\\" xmlns=\\\"http://schemas.microsoft.com/windows/2004/02/mit/task\\\">\\r\\n  <RegistrationInfo>\\r\\n    <Date>2026-07-07T23:33:53</Date>\\r\\n    <Author>hostname\\\\username</Author>\\r\\n    <URI>\\\\T1053_005_OnStartup</URI>\\r\\n  </RegistrationInfo>\\r\\n  <Triggers>\\r\\n    <BootTrigger>\\r\\n      <StartBoundary>2026-07-07T23:33:00</StartBoundary>\\r\\n      <Enabled>true</Enabled>\\r\\n    </BootTrigger>\\r\\n  </Triggers>\\r\\n  <Settings>\\r\\n    <MultipleInstancesPolicy>IgnoreNew</MultipleInstancesPolicy>\\r\\n    <DisallowStartIfOnBatteries>true</DisallowStartIfOnBatteries>\\r\\n    <StopIfGoingOnBatteries>true</StopIfGoingOnBatteries>\\r\\n    <AllowHardTerminate>true</AllowHardTerminate>\\r\\n    <StartWhenAvailable>false</StartWhenAvailable>\\r\\n    <RunOnlyIfNetworkAvailable>false</RunOnlyIfNetworkAvailable>\\r\\n    <IdleSettings>\\r\\n      <Duration>PT10M</Duration>\\r\\n      <WaitTimeout>PT1H</WaitTimeout>\\r\\n      <StopOnIdleEnd>true</StopOnIdleEnd>\\r\\n      <RestartOnIdle>false</RestartOnIdle>\\r\\n    </IdleSettings>\\r\\n    <AllowStartOnDemand>true</AllowStartOnDemand>\\r\\n    <Enabled>true</Enabled>\\r\\n    <Hidden>false</Hidden>\\r\\n    <RunOnlyIfIdle>false</RunOnlyIfIdle>\\r\\n    <WakeToRun>false</WakeToRun>\\r\\n    <ExecutionTimeLimit>PT72H</ExecutionTimeLimit>\\r\\n    <Priority>7</Priority>\\r\\n  </Settings>\\r\\n  <Actions Context=\\\"Author\\\">\\r\\n    <Exec>\\r\\n      <Command>cmd.exe</Command>\\r\\n      <Arguments>/c calc.exe</Arguments>\\r\\n    </Exec>\\r\\n  </Actions>\\r\\n  <Principals>\\r\\n    <Principal id=\\\"Author\\\">\\r\\n      <UserId>S-1-5-18</UserId>\\r\\n      <RunLevel>LeastPrivilege</RunLevel>\\r\\n    </Principal>\\r\\n  </Principals>\\r\\n</Task>\\r\\n\\r\\nOther Information:\\r\\n\\tProcessCreationTime: \\t\\t9007199254743556\\r\\n\\tClientProcessId: \\t\\t\\t5376\\r\\n\\tParentProcessId: \\t\\t\\t6080\\r\\n\\tFQDN: \\t\\t0\\r\\n\\t\\\"\"},\"eventdata\":{\"subjectUserSid\":\"S-1-5-21-1936990016-347023981-4024709258-1001\",\"subjectUserName\":\"username\",\"subjectDomainName\":\"hostname\",\"subjectLogonId\":\"0x32b7647\",\"taskName\":\"\\\\\\\\T1053_005_OnStartup\",\"taskContent\":\"&lt;?xml version=\\\\\\\"1.0\\\\\\\" encoding=\\\\\\\"UTF-16\\\\\\\"?&gt;  &lt;Task version=\\\\\\\"1.2\\\\\\\" xmlns=\\\\\\\"http://schemas.microsoft.com/windows/2004/02/mit/task\\\\\\\"&gt;    &lt;RegistrationInfo&gt;      &lt;Date&gt;2026-07-07T23:33:53&lt;/Date&gt;      &lt;Author&gt;hostname\\\\\\\\username&lt;/Author&gt;      &lt;URI&gt;\\\\\\\\T1053_005_OnStartup&lt;/URI&gt;    &lt;/RegistrationInfo&gt;    &lt;Triggers&gt;      &lt;BootTrigger&gt;        &lt;StartBoundary&gt;2026-07-07T23:33:00&lt;/StartBoundary&gt;        &lt;Enabled&gt;true&lt;/Enabled&gt;      &lt;/BootTrigger&gt;    &lt;/Triggers&gt;    &lt;Settings&gt;      &lt;MultipleInstancesPolicy&gt;IgnoreNew&lt;/MultipleInstancesPolicy&gt;      &lt;DisallowStartIfOnBatteries&gt;true&lt;/DisallowStartIfOnBatteries&gt;      &lt;StopIfGoingOnBatteries&gt;true&lt;/StopIfGoingOnBatteries&gt;      &lt;AllowHardTerminate&gt;true&lt;/AllowHardTerminate&gt;      &lt;StartWhenAvailable&gt;false&lt;/StartWhenAvailable&gt;      &lt;RunOnlyIfNetworkAvailable&gt;false&lt;/RunOnlyIfNetworkAvailable&gt;      &lt;IdleSettings&gt;        &lt;Duration&gt;PT10M&lt;/Duration&gt;        &lt;WaitTimeout&gt;PT1H&lt;/WaitTimeout&gt;        &lt;StopOnIdleEnd&gt;true&lt;/StopOnIdleEnd&gt;        &lt;RestartOnIdle&gt;false&lt;/RestartOnIdle&gt;      &lt;/IdleSettings&gt;      &lt;AllowStartOnDemand&gt;true&lt;/AllowStartOnDemand&gt;      &lt;Enabled&gt;true&lt;/Enabled&gt;      &lt;Hidden&gt;false&lt;/Hidden&gt;      &lt;RunOnlyIfIdle&gt;false&lt;/RunOnlyIfIdle&gt;      &lt;WakeToRun&gt;false&lt;/WakeToRun&gt;      &lt;ExecutionTimeLimit&gt;PT72H&lt;/ExecutionTimeLimit&gt;      &lt;Priority&gt;7&lt;/Priority&gt;    &lt;/Settings&gt;    &lt;Actions Context=\\\\\\\"Author\\\\\\\"&gt;      &lt;Exec&gt;        &lt;Command&gt;cmd.exe&lt;/Command&gt;        &lt;Arguments&gt;/c calc.exe&lt;/Arguments&gt;      &lt;/Exec&gt;    &lt;/Actions&gt;    &lt;Principals&gt;      &lt;Principal id=\\\\\\\"Author\\\\\\\"&gt;        &lt;UserId&gt;S-1-5-18&lt;/UserId&gt;        &lt;RunLevel&gt;LeastPrivilege&lt;/RunLevel&gt;      &lt;/Principal&gt;    &lt;/Principals&gt;  &lt;/Task&gt;\",\"clientProcessStartKey\":\"9007199254743556\",\"clientProcessId\":\"5376\",\"parentProcessId\":\"6080\",\"rpcCallClientLocality\":\"0\",\"fQDN\":\"hostname\"}}}",
  "input": {
    "type": "log"
  },
  "@timestamp": "2026-07-08T06:33:54.332Z",
  "location": "EventChannel",
  "id": "1783492434.1598915",
  "timestamp": "2026-07-08T06:33:54.332+0000"
}