# sigma/persistence/T1053.005_scheduled_task/T1053005_task_hidden.yml
title: Hidden Scheduled Task Created
id: b2c3d4e5-f6a7-8901-bcde-222222222222
status: experimental
description: >
    Detects scheduled tasks created with the Hidden property set
    to true. A hidden task doesn't show up in the Task Scheduler
    GUI, so attackers use it to keep persistence out of sight from
    anyone checking by hand.
references:
    - https://github.com/jeremiah-onyema/windows-detection-research
    - https://attack.mitre.org/techniques/T1053/005/
author: Jeremiah
date: 2026/03/26
tags:
    - attack.persistence
    - attack.defense_evasion
    - attack.t1053.005
logsource:
    product: windows
    service: security
    definition: 'Requirements: Windows Event ID 4698'
detection:
    selection:
        EventID: 4698
    selection_hidden:
        TaskContent|contains:
            - '<Hidden>true</Hidden>'
    condition: selection and selection_hidden
falsepositives:
    - Some legitimate software creates hidden
      scheduled tasks for background operations
level: high
