# evidence/T1547.001/secedit_services_exe_run_key.json
{
  "agent": {
    "ip": "x.x.x.x",
    "name": "windows-detection-vm",
    "id": "001"
  },
  "manager": {
    "name": "wazuh-manager"
  },
  "data": {
    "win": {
      "eventdata": {
        "image": "C:\\\\WINDOWS\\\\system32\\\\services.exe",
        "targetObject": "HKLM\\\\SOFTWARE\\\\Microsoft\\\\Windows\\\\CurrentVersion\\\\Run\\\\calc",
        "processGuid": "{f5cafcf2-e7a0-6a52-0b00-000000002300}",
        "processId": "976",
        "utcTime": "2026-07-11 16:04:53.324",
        "ruleName": "technique_id=T1547.001,technique_name=Registry Run Keys / Start Folder",
        "details": "calc.exe",
        "eventType": "SetValue",
        "user": "NT AUTHORITY\\\\SYSTEM"
      },
      "system": {
        "eventID": "13",
        "keywords": "0x8000000000000000",
        "providerGuid": "{5770385f-c22a-43e0-bf4c-06f5698ffbd9}",
        "level": "4",
        "channel": "Microsoft-Windows-Sysmon/Operational",
        "opcode": "0",
        "message": "\"Registry value set:\r\nRuleName: technique_id=T1547.001,technique_name=Registry Run Keys / Start Folder\r\nEventType: SetValue\r\nUtcTime: 2026-07-11 16:04:53.324\r\nProcessGuid: {f5cafcf2-e7a0-6a52-0b00-000000002300}\r\nProcessId: 976\r\nImage: C:\\WINDOWS\\system32\\services.exe\r\nTargetObject: HKLM\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run\\calc\r\nDetails: calc.exe\r\nUser: NT AUTHORITY\\SYSTEM\"",
        "version": "2",
        "systemTime": "2026-07-11T16:04:53.3268052Z",
        "eventRecordID": "2821805",
        "threadID": "4648",
        "computer": "hostname",
        "task": "13",
        "processID": "3788",
        "severityValue": "INFORMATION",
        "providerName": "Microsoft-Windows-Sysmon"
      }
    }
  },
  "rule": {
    "firedtimes": 1,
    "mail": true,
    "level": 15,
    "description": "High Risk: Unexpected process wrote executable to Run key: C:\\\\WINDOWS\\\\system32\\\\services.exe",
    "groups": [
      "windows",
      "persistence",
      "sigma_coverage"
    ],
    "mitre": {
      "technique": [
        "Registry Run Keys / Startup Folder"
      ],
      "id": [
        "T1547.001"
      ],
      "tactic": [
        "Persistence",
        "Privilege Escalation"
      ]
    },
    "id": "100013"
  },
  "decoder": {
    "name": "windows_eventchannel"
  },
  "full_log": "{\"win\":{\"system\":{\"providerName\":\"Microsoft-Windows-Sysmon\",\"providerGuid\":\"{5770385f-c22a-43e0-bf4c-06f5698ffbd9}\",\"eventID\":\"13\",\"version\":\"2\",\"level\":\"4\",\"task\":\"13\",\"opcode\":\"0\",\"keywords\":\"0x8000000000000000\",\"systemTime\":\"2026-07-11T16:04:53.3268052Z\",\"eventRecordID\":\"2821805\",\"processID\":\"3788\",\"threadID\":\"4648\",\"channel\":\"Microsoft-Windows-Sysmon/Operational\",\"computer\":\"hostname\",\"severityValue\":\"INFORMATION\",\"message\":\"\\\"Registry value set:\\r\\nRuleName: technique_id=T1547.001,technique_name=Registry Run Keys / Start Folder\\r\\nEventType: SetValue\\r\\nUtcTime: 2026-07-11 16:04:53.324\\r\\nProcessGuid: {f5cafcf2-e7a0-6a52-0b00-000000002300}\\r\\nProcessId: 976\\r\\nImage: C:\\\\WINDOWS\\\\system32\\\\services.exe\\r\\nTargetObject: HKLM\\\\SOFTWARE\\\\Microsoft\\\\Windows\\\\CurrentVersion\\\\Run\\\\calc\\r\\nDetails: calc.exe\\r\\nUser: NT AUTHORITY\\\\SYSTEM\\\"\"},\"eventdata\":{\"ruleName\":\"technique_id=T1547.001,technique_name=Registry Run Keys / Start Folder\",\"eventType\":\"SetValue\",\"utcTime\":\"2026-07-11 16:04:53.324\",\"processGuid\":\"{f5cafcf2-e7a0-6a52-0b00-000000002300}\",\"processId\":\"976\",\"image\":\"C:\\\\\\\\WINDOWS\\\\\\\\system32\\\\\\\\services.exe\",\"targetObject\":\"HKLM\\\\\\\\SOFTWARE\\\\\\\\Microsoft\\\\\\\\Windows\\\\\\\\CurrentVersion\\\\\\\\Run\\\\\\\\calc\",\"details\":\"calc.exe\",\"user\":\"NT AUTHORITY\\\\\\\\SYSTEM\"}}}",
  "input": {
    "type": "log"
  },
  "@timestamp": "2026-07-11T16:04:53.794Z",
  "location": "EventChannel",
  "id": "1783785893.370941",
  "timestamp": "2026-07-11T16:04:53.794+0000"
}