# evidence/T1547.009/startup_lnk_programdata.json
{
  "agent": {
    "ip": "x.x.x.x",
    "name": "windows-detection-vm",
    "id": "001"
  },
  "manager": {
    "name": "wazuh-manager"
  },
  "data": {
    "win": {
      "eventdata": {
        "image": "C:\\\\WINDOWS\\\\System32\\\\WindowsPowerShell\\\\v1.0\\\\powershell.exe",
        "processGuid": "{f5cafcf2-9634-6a4e-e701-000000002100}",
        "processId": "4944",
        "utcTime": "2026-07-08 18:25:56.780",
        "targetFilename": "C:\\\\ProgramData\\\\Microsoft\\\\Windows\\\\Start Menu\\\\Programs\\\\Startup\\\\T1547.009.lnk",
        "ruleName": "technique_id=T1187,technique_name=Forced Authentication",
        "creationUtcTime": "2026-07-08 18:25:56.780",
        "user": "hostname\\\\username"
      },
      "system": {
        "eventID": "11",
        "keywords": "0x8000000000000000",
        "providerGuid": "{5770385f-c22a-43e0-bf4c-06f5698ffbd9}",
        "level": "4",
        "channel": "Microsoft-Windows-Sysmon/Operational",
        "opcode": "0",
        "message": "\"File created:\r\nRuleName: technique_id=T1187,technique_name=Forced Authentication\r\nUtcTime: 2026-07-08 18:25:56.780\r\nProcessGuid: {f5cafcf2-9634-6a4e-e701-000000002100}\r\nProcessId: 4944\r\nImage: C:\\WINDOWS\\System32\\WindowsPowerShell\\v1.0\\powershell.exe\r\nTargetFilename: C:\\ProgramData\\Microsoft\\Windows\\Start Menu\\Programs\\Startup\\T1547.009.lnk\r\nCreationUtcTime: 2026-07-08 18:25:56.780\r\nUser: hostname\\username\"",
        "version": "2",
        "systemTime": "2026-07-08T18:25:56.7814505Z",
        "eventRecordID": "2802673",
        "threadID": "4456",
        "computer": "hostname",
        "task": "11",
        "processID": "3596",
        "severityValue": "INFORMATION",
        "providerName": "Microsoft-Windows-Sysmon"
      }
    }
  },
  "rule": {
    "firedtimes": 2,
    "mail": true,
    "level": 15,
    "description": "High Risk: Suspicious process created startup shortcut: C:\\\\WINDOWS\\\\System32\\\\WindowsPowerShell\\\\v1.0\\\\powershell.exe wrote C:\\\\ProgramData\\\\Microsoft\\\\Windows\\\\Start Menu\\\\Programs\\\\Startup\\\\T1547.009.lnk",
    "groups": [
      "windows",
      "persistence",
      "sigma_coverage"
    ],
    "mitre": {
      "technique": [
        "Shortcut Modification"
      ],
      "id": [
        "T1547.009"
      ],
      "tactic": [
        "Persistence",
        "Privilege Escalation"
      ]
    },
    "id": "100016"
  },
  "decoder": {
    "name": "windows_eventchannel"
  },
  "input": {
    "type": "log"
  },
  "@timestamp": "2026-07-08T18:25:57.619Z",
  "location": "EventChannel",
  "id": "1783535157.3435743",
  "timestamp": "2026-07-08T18:25:57.619+0000"
}